# wisetime/crontab
1 * * * * python /app/manage.py runcrons > /home/cronjob.log